# packages/service/minecraft-be-server.service
[Unit]
Description=Minecraft Server BE with MCSO
PartOf=minecraft-master-session.service
After=Network-online.target minecraft-master-session.service

[Service]
Type=oneshot
ExecStart=/opt/mcso/start_mcso.sh -b
ExecStop=/opt/mcso/stop_mcso.sh -b
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target minecraft-master-session.service